# supervisor-config.conf
[program:whatsmark-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/whatsmark/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/whatsmark/storage/logs/queue.log
stopwaitsecs=3600

# Save this as: /etc/supervisor/conf.d/whatsmark-queue.conf
# Then run:
# supervisorctl reread
# supervisorctl update
# supervisorctl start whatsmark-queue:*